# systemd/deepswing.service
[Unit]
Description=DeepSwing AI Swing Trading Simulator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alexander
WorkingDirectory=/home/alexander/Documents/DeepSwing
ExecStart=/home/alexander/Documents/DeepSwing/venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits appropriate for Pi 5
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
